# bsads_backend_and_fast_api/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:postgresql]
user=postgres
command=/app/start_postgres.sh
autostart=true
autorestart=true
priority=1
stdout_logfile=/var/log/supervisor/postgresql.log
stderr_logfile=/var/log/supervisor/postgresql_err.log

[program:fastapi]
command=/app/start_production.sh
autostart=true
autorestart=true
priority=2
stdout_logfile=/var/log/supervisor/fastapi.log
stderr_logfile=/var/log/supervisor/fastapi_err.log
# Wait for PostgreSQL to be ready before starting FastAPI
startsecs=10
